# scripts/companion/companion-model-proxy.service
# /etc/systemd/system/companion-model-proxy.service
# Install:
#   sudo cp scripts/companion/companion-model-proxy.service /etc/systemd/system/
#   sudo systemctl daemon-reload && sudo systemctl enable --now companion-model-proxy
[Unit]
Description=Crow Companion — model-routing proxy (fast 4B -> 35B)
After=network-online.target crow-gateway.service
Wants=network-online.target

[Service]
Type=simple
User=kh0pp
WorkingDirectory=/home/kh0pp/crow
Environment=PATH=/home/kh0pp/.nvm/versions/node/v20.20.2/bin:/usr/local/bin:/usr/bin:/bin
Environment=CROW_DB_PATH=/home/kh0pp/.crow/data/crow.db
Environment=COMPANION_PROXY_PORT=11435
Environment=COMPANION_FAST_MODEL=crow-voice/qwen3.5-4b
Environment=COMPANION_ESCALATION_MODEL=crow-chat/qwen3.6-35b-a3b
ExecStart=/home/kh0pp/.nvm/versions/node/v20.20.2/bin/node /home/kh0pp/crow/scripts/companion/model-proxy.mjs
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
