# scripts/pi-bots/systemd/pibot-bridge@.service
[Unit]
Description=Crow Bot Builder — Gmail bridge tick [%i]
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=kh0pp
WorkingDirectory=/home/kh0pp/crow
EnvironmentFile=/etc/crow/pibot-%i.env
ExecStart=/usr/bin/env node /home/kh0pp/crow/scripts/pi-bots/bridge_tick.mjs
TimeoutStartSec=600
StandardOutput=journal
StandardError=journal
